# payload/Library/LaunchAgents/com.grahamgilbert.chrome_update_notifier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.grahamgilbert.chrome_update_notifier</string>
    <key>Program</key>
    <string>/Library/Management/chrome_update_notifier.py</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>300</integer>
  </dict>
</plist>
